# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:celeryd]
command=/home/sardor/moffice/billing/envbilling/bin/celery worker
autostart=true
autorestart=unexpected
stderr_logfile=celery_error.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10

